# templates/default/upstart.conf.erb
description "<%= @app_name %>"

start on runlevel [234]
stop on runlevel [!234]

setuid <%= @user %>
setgid <%= @group %>

env APP_DIR=<%= @deploy_to %>

respawn

script
exec /bin/bash -l <<'EOT'
  source <%= @deploy_to %>/environment
  cd <%= @deploy_to %>/current
  exec bundle exec <%= @exec_command %> >>$APP_DIR/shared/log/upstart.log 2>&1
EOT
end script
