# templates/default/upstart.unicorn.conf.erb
description "<%= @app_name %>"

start on runlevel [234]
stop on runlevel [!234]

setuid <%= @user %>
setgid <%= @group %>

env APP_DIR=<%= @deploy_to %>

respawn

script
exec /bin/bash -l <<'EOT'
  while true; do
    if [ ! -f $APP_DIR/shared/pids/unicorn.<%= @app_name %>.pid ]; then
      # if pid file does not exist, start unicorn
      source <%= @deploy_to %>/environment
      cd $APP_DIR/current
      bundle exec <%= @exec_command %> >>$APP_DIR/shared/log/upstart.log 2>&1
    else
      # if pid file does not exist, just hang around waiting to get a lock
      flock -x 0 < $APP_DIR/shared/pids/unicorn.<%= @app_name %>.pid.lock
      # if we got here, master exited (died or restarted)
      sleep 2
    fi
  done
EOT
end script
post-stop script
  if [ -f $APP_DIR/shared/pids/unicorn.<%= @app_name %>.pid ]; then
    kill $(cat $APP_DIR/shared/pids/unicorn.<%= @app_name %>.pid)
  fi
end script
